# gpsd-nmea-bridge.socket
[Unit]
Description=Apache Server Socket

[Socket]
ListenStream=10110

[Install]
WantedBy=sockets.target